# deploy/gunicorn.socket
[Unit]
Description=Phonix Gunicorn Socket

[Socket]
ListenStream=/run/phonix/gunicorn.sock
SocketUser=www-data

[Install]
WantedBy=sockets.target
